# rootfs/etc/local.d/claude-setup.start
#!/bin/sh
# clawdOS early setup
# Ensure config directory exists
mkdir -p /etc/clawdos
# Ensure correct permissions on env file
[ -f /etc/clawdos/env ] && chmod 0600 /etc/clawdos/env
